# docker/s6/services/haproxy-api.svc/run
#!/bin/sh

# stderr -> stdout
exec 2>&1

cd /haproxy-api

exec ./haproxy-api -f haproxy-api.toml
