# docker/s6/services/haproxy.svc/run
#!/bin/sh

# stderr -> stdout
exec 2>&1

cd /

exec haproxy -f /etc/haproxy.cfg
